# root/etc/s6-overlay/s6-rc.d/svc-mod-universal-openssh/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc \
      sshd -D -e