# root/etc/s6-overlay/s6-rc.d/init-mod-universal-openssh-install/run
#!/usr/bin/with-contenv bash

echo "**** Setting up openssh ****"
ssh-keygen -A
